# crontab
0 0 * * * /usr/local/bin/python /home/weatheralmanac/source/main.py
