# crontab
0 0 * * * docker exec weatheralmanac python3 main.py > ../log/weatheralmanac.log 2>&1
